# deploy/flibusta-bot.service
[Unit]
Description=Telegram Flibusta Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/main.py
Restart=always
RestartSec=10
RuntimeMaxSec=24h

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=__APP_DIR__

[Install]
WantedBy=multi-user.target
